# deploy/pmrec.service
[Unit]
Description=Polymarket research data recorder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust user and paths for your server.
User=pmrec
WorkingDirectory=/opt/polymarket-recorder
Environment=PMREC_CONFIG=/opt/polymarket-recorder/config.yaml
# Use the venv's python.
ExecStart=/opt/polymarket-recorder/.venv/bin/python -m pmrec.main
Restart=always
RestartSec=3
# Capture is the source of truth; restart aggressively and never give up.
StartLimitIntervalSec=0
# Give it room; uploads + buffers are memory-light but be generous.
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
